# trading_bot.service
[Unit]
Description=Gunicorn instance to serve AI Trading Bot
After=network.target

[Service]
# IMPORTANT: Ensure /home/ubuntu/Ai-Bot-Trader is owned by ubuntu:www-data
# Run: sudo chown -R ubuntu:www-data /home/ubuntu/Ai-Bot-Trader
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/Ai-Bot-Trader
Environment="PATH=/home/ubuntu/Ai-Bot-Trader/venv/bin"
ExecStart=/home/ubuntu/Ai-Bot-Trader/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 run:app

[Install]
WantedBy=multi-user.target
